# board/runtipi/common/rootfs-overlay/etc/systemd/system/runtipios-hotspot.service
[Unit]
Description=RuntipiOS WiFi Hotspot
After=network.target
Conflicts=wpa_supplicant.service

[Service]
Type=forking
ExecStart=/usr/local/bin/runtipios-hotspot start
ExecStop=/usr/local/bin/runtipios-hotspot stop
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
