# runtipios/scripts/runtipi-wifi.service
[Unit]
Description=Runtipi Wi-Fi Detection and Setup Service
After=network-online.target
Wants=network-online.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/runtipi-hotspot/scripts/detect-network.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
