# board/runtipi/rootfs-overlay/etc/systemd/system/runtipi-motd.service
[Unit]
Description=RuntipiOS MOTD Generator
After=runtipi-install.service
Wants=runtipi-install.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/runtipi-generate-motd.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
