# board/runtipi/rootfs-overlay/etc/systemd/system/runtipi-network.service
[Unit]
Description=RuntipiOS Network Detection and Setup
After=network-pre.target
Before=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/runtipi-network-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
