# board/runtipi/rootfs-overlay/etc/systemd/system/runtipi-captive.service
[Unit]
Description=RuntipiOS Captive Portal (WiFi Setup)
After=systemd-networkd.service
Before=runtipi-install.service
BindsTo=systemd-networkd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/runtipi-captive-portal.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target
