# platforms/common/overlays/etc/systemd/system/runtipi-install.service
[Unit]
Description=Runtipi Auto-Installation Service
After=network-online.target runtipios-setup.service
Wants=network-online.target

[Service]
Type=oneshot
User=runtipi
# Wait for internet connectivity
ExecStartPre=/bin/bash -c 'until ping -c1 8.8.8.8 &>/dev/null; do sleep 5; done'
ExecStart=/usr/local/bin/runtipi-install
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
