# platforms/common/overlays/etc/systemd/system/runtipios-setup.service
[Unit]
Description=RuntipiOS First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/etc/runtipios-configured

[Service]
Type=oneshot
ExecStart=/usr/local/bin/runtipios-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
